# contrib/init/bellsd.conf
description "TenebrixCoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env tenebrixcoind_BIN="/usr/bin/tenebrixcoind"
env tenebrixcoind_USER="tenebrixcoin"
env tenebrixcoind_GROUP="tenebrixcoin"
env tenebrixcoind_PIDDIR="/var/run/tenebrixcoind"
# upstart can't handle variables constructed with other variables
env tenebrixcoind_PIDFILE="/var/run/tenebrixcoind/tenebrixcoind.pid"
env tenebrixcoind_CONFIGFILE="/etc/tenebrixcoin/tenebrixcoin.conf"
env tenebrixcoind_DATADIR="/var/lib/tenebrixcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # tenebrixcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$tenebrixcoind_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run tenebrixcoind."
        echo "The setting must appear in $tenebrixcoind_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"TenebrixCoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$tenebrixcoind_PIDDIR"
    chmod 0755 "$tenebrixcoind_PIDDIR"
    chown $tenebrixcoind_USER:$tenebrixcoind_GROUP "$tenebrixcoind_PIDDIR"
    chown $tenebrixcoind_USER:$tenebrixcoind_GROUP "$tenebrixcoind_CONFIGFILE"
    chmod 0660 "$tenebrixcoind_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$tenebrixcoind_PIDFILE" \
    --chuid $tenebrixcoind_USER:$tenebrixcoind_GROUP \
    --exec "$tenebrixcoind_BIN" \
    -- \
    -pid="$tenebrixcoind_PIDFILE" \
    -conf="$tenebrixcoind_CONFIGFILE" \
    -datadir="$tenebrixcoind_DATADIR" \
    -disablewallet \
    -daemon

